# bin/mcsds.service
[Unit]
Description=MCS Data Switch
RequiresMountsFor=/etc/cpdisplay/config
After=network.target

[Service]
ExecStart=/opt/cpdisplay/bin/mcsdataswitch -c /etc/cpdisplay/config/vscl/cpd-settings.mcsdataswitch.ini
Restart=on-abort

[Install]
WantedBy=multi-user.target
